# launchd/com.example.sesame-remo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.sesame-remo</string>
  <key>ProgramArguments</key>
  <array>
    <string>/absolute/path/to/.venv/bin/python</string>
    <string>-m</string>
    <string>sesame_remo</string>
    <string>touch-pro-trigger</string>
    <string>--config</string>
    <string>/absolute/path/to/config.toml</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/absolute/path/to/project</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/tmp/sesame-remo.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/sesame-remo.err.log</string>
</dict>
</plist>
